# opennms-base-assembly/src/main/resources/contrib/failover/etc/sync-on.crontab
#Crontab entry for automatic synchronization of OpenNMS
MAILTO=root
SYNC_USER=opennms
0 0 * * *	/opt/opennms/contrib/failover/scripts/sync.sh
