# debian/opennms-server.prerm
#! /bin/sh
# prerm script for opennms-server
#
# $Id$

set -e

case "$1" in
    remove|upgrade|deconfigure)
	invoke-rc.d opennms stop || :
	;;
    failed-upgrade)
	;;
    *)
	echo "prerm called with unknown argument \`$1'" >&2
	exit 1
	;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


