# debian/postinst
#! /bin/bash
# postinst script for opennms
#
# see: dh_installdeb(1)

set -e

case "$1" in
    configure)
	if [ -x /sbin/ldconfig ] ; then
		ldconfig
	fi

	if [ -x "/etc/init.d/opennms" ]; then
		update-rc.d opennms defaults >/dev/null
	fi

	# FIXME: add parsing of /etc/postgresql/* stuff
	/etc/init.d/postgresql restart

	echo -e "- running the post-install..."

	/usr/share/opennms/bin/install.pl -r -q -l /usr/lib/postgresql/lib/opennms -u opennms -p opennms -U postgres -d opennms /etc/opennms/create.sql

	echo "If this is a first-time install, please modify the /etc/opennms"
	echo "files before starting OpenNMS for the first time."

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


