# debian/postrm
#! /bin/sh
# postrm script for opennms
#
# see: dh_installdeb(1)

set -e

case "$1" in
    purge|remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
	if [ -x /sbin/ldconfig ] ; then
		ldconfig
	fi

        if [ "$1" = "purge" ]; then
            update-rc.d opennms remove >/dev/null

            echo -e "- removing working directories... \c"
            rm -rf /var/log/opennms /var/lib/opennms >/dev/null || :
            echo "done"

            if [ -x /usr/lib/postgresql/bin/dropdb ]; then
                echo -n "removing the opennms database tables..."
                /usr/lib/postgresql/bin/dropdb "opennms" >/dev/null || :
                echo " done."
            fi
        fi
        ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 0

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#


